# urlwatch/rootfs/etc/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true

[program:deluged]
user=urlwatch
environment=TERM="xterm"
command=/usr/bin/urlwatch