# urlwatch/rootfs/etc/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true

[program:urlwatch]
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
user=root
environment=TERM="xterm"
command=gosu ${UID}:${GID} sh -c "/usr/bin/urlwatch1"