# urlwatch/rootfs/etc/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true

[program:urlwatch]
user=urlwatch
environment=TERM="xterm"
command=/usr/bin/urlwatch